# debian/etc/cron.d/gvm
0  2    * * *   root    flock -xn /opt/gvm/update/update.lock -c /opt/gvm/update/update.sh > /dev/null 2> /opt/gvm/update/update.log